# arch-niri-noctalia-v6.2/files/home/.config/systemd/user/noctalia.service
[Unit]
Description=Noctalia Shell
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/qs -c noctalia-shell
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
